# deploy/crm-backend.service
[Unit]
Description=CRM RUTEST FastAPI backend
After=network.target

[Service]
Type=simple
User=crm
Group=crm
WorkingDirectory=/opt/crm/backend
EnvironmentFile=/opt/crm/.env
ExecStart=/opt/crm/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=append:/var/log/crm/backend.log
StandardError=append:/var/log/crm/backend.err.log

# Безопасность
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
